# deploy/nora-demo-traffic.service
[Unit]
Description=NORA Demo Traffic Simulator
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/opt/nora/demo-traffic.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
